# systemd/prepare-immich-backup.service
[Unit]
Description=QuesadaLab - Prepare Immich Backup Metadata
After=docker.service network-online.target srv-immich\x2ddata.mount
Requires=docker.service
RequiresMountsFor=/srv/immich-data
Wants=network-online.target

[Service]
Type=oneshot
User=root
Group=root
Environment=IMMICH_BACKUP_ROOT=/opt/quesadalab/backups/immich
Environment=IMMICH_PREPARED_RETENTION=3
ExecStart=/opt/quesadalab/scripts/prepare-immich-backup.sh
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
TimeoutStartSec=2h
PrivateTmp=yes
ProtectHome=yes
ProtectSystem=full
NoNewPrivileges=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
