# systemd/pull-immich-backups.timer
[Unit]
Description=QuesadaLab - Daily Immich USB Pull Timer

[Timer]
OnCalendar=*-*-* 01:15:00
Persistent=true
AccuracySec=5min
Unit=pull-immich-backups.service

[Install]
WantedBy=timers.target
