# systemd/pull-jellyfin-backups.service
[Unit]
Description=QuesadaLab - Pull Jellyfin Backups to USB
After=network-online.target mnt-quesadalab\x2dbackup.mount
Wants=network-online.target
RequiresMountsFor=/mnt/quesadalab-backup
ConditionPathIsMountPoint=/mnt/quesadalab-backup

[Service]
Type=oneshot
User=root
Group=root
Environment=JELLYFIN_USB_RETENTION=7
ExecStart=/usr/local/sbin/pull-jellyfin-backups.sh
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
TimeoutStartSec=2h
PrivateTmp=yes
ProtectHome=read-only
ProtectSystem=full
NoNewPrivileges=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
